# MacPatch PKG/Updater/PKG_SRC/Library/LaunchDaemons/gov.llnl.mp.agentUpdater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>gov.llnl.mp.agentUpdater</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/MacPatch/Updater/MPAgentUp2Date</string>
		<string>-u</string>
	</array>
	<key>ServiceDescription</key>
	<string>MacPatch Apple Patches Downloader</string>
	<key>StandardErrorPath</key>
	<string>/Library/MacPatch/Updater/Logs/stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Library/MacPatch/Updater/Logs/stdout.log</string>
	<key>StartInterval</key>
	<integer>5400</integer>
	<key>onDemand</key>
	<true/>
	<key>QueueDirectories</key>
	<array>
		<string>/Users/Shared/.mpUpdate</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
