# Packages/Client/Scripts/postinstall
#!/bin/sh -x

# ***************************************************************************************
# PostFlight Script for MacPatch Client Installer
# ver 3.0.0
# Created By Charles Heizer @ LLNL
#
# ***************************************************************************************

MP_ROOT="/Library/MacPatch"
MP_ROOT_CLIENT="${MP_ROOT}/Client"

# Script Parent Dir Path
SRC="$(dirname "$0")"

# Current logged in User ID
# consoleUserID at login window is 0
consoleUserID=$(stat -f%u /dev/console)

setOwner() {
	if [ -f "$1" ]; then
		/usr/sbin/chown -R root:wheel "$1"
	fi
}

setOwnerAndPerms() {
	if [ -f "$1" ]; then
		/usr/sbin/chown root:wheel "$1"
		/bin/chmod 0644 "$1"
	fi
}

# ******************************************************************
# Copy in client version info
if [ -e "${SRC}/.mpVersion.plist" ]; then
	cp "${SRC}/.mpVersion.plist" "${MP_ROOT_CLIENT}/.mpVersion.plist"
fi

# ******************************************************************
# Set owner & permissions
setOwnerAndPerms /Library/LaunchDaemons/gov.llnl.mp.helper.plist
setOwnerAndPerms /Library/LaunchAgents/gov.llnl.mp.status.plist
setOwner /Library/PrivilegedHelperTools/gov.llnl.mp.helper


# ******************************************************************
# Launch the Helper service  ***************************************
#
/bin/launchctl load -w /Library/LaunchDaemons/gov.llnl.mp.helper.plist
/bin/sleep 2

# ******************************************************************
# Start LaunchAgent Apps *******************************************
#
if [ "$consoleUserID" != "0" ]; then
  /bin/launchctl bootstrap gui/$consoleUserID /Library/LaunchAgents/gov.llnl.mp.status.plist
fi

# Kill the ClientStatus app, launchctl does not restart the app only the service
/usr/bin/killall MPClientStatus

exit 0;
