# Source/MacPatch/Packages/Base/Scripts/postinstall
#!/bin/sh -x

# ***************************************************************************************
# PostFlight Script for MacPatch Client Installer
# ver 3.0.0
# Created By Charles Heizer @ LLNL
#
# ***************************************************************************************

MP_ROOT="/Library/MacPatch"
MP_ROOT_CLIENT="/Library/MacPatch/Client"

# Script Parent Dir Path
SRC=$(cd $(dirname "$0"); pwd)

# System Version Info
sysVersion=$(uname -r)
sysMajorVersion=${sysVersion%%.*}
tempMinorVersion=${sysVersion#*.}
sysMinorVersion=${tempMinorVersion%%.*}

# Current Logged in user info
usrID=$(stat -f%Su /dev/console)
if [ $? -ne 0 ]; then
usrID="NULL"
fi

function setOwner
{
	if [ -f "$1" ]; then
		/usr/sbin/chown -R root:wheel $1
	fi
}

function setOwnerAndPerms
{
	if [ -f "$1" ]; then
		/usr/sbin/chown root:wheel $1
		/bin/chmod 0644 $1
	fi
}

# ******************************************************************
# Copy in client version info
if [ -e "${SRC}/.mpVersion.plist" ]; then
	cp "${SRC}/.mpVersion.plist" "/Library/MacPatch/Client/.mpVersion.plist"
fi

# ******************************************************************
# Set owner & permissions
setOwnerAndPerms /Library/LaunchDaemons/gov.llnl.mp.helper.plist
setOwnerAndPerms /Library/LaunchAgents/gov.llnl.mp.status.plist
setOwner /Library/PrivilegedHelperTools/gov.llnl.mp.helper


# ******************************************************************
# Launch the Helper service  ***************************************
#
/bin/launchctl load -w /Library/LaunchDaemons/gov.llnl.mp.helper.plist
/bin/sleep 2

# ******************************************************************
# Start LaunchAgent Apps *******************************************
#
if [ "$usrID" != "null" ]; then
	su -l $usrID -c 'launchctl load -wF /Library/LaunchAgents/gov.llnl.mp.status.plist'
fi

# Kill the ClientStatus app, launchctl does not restart the app only the service
/usr/bin/killall MPClientStatus

exit 0;
