# Source/Server/conf/systemd/MPRsyncServer3.socket
[Unit]
Description=MacPatch Rsync Server Socket
Conflicts=MPRsyncServer3.service

[Socket]
ListenStream=873
Accept=yes

[Install]
WantedBy=sockets.target 
