# Source/Client/planb/gov.llnl.mp.planb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>gov.llnl.mp.planb</string>
	<key>LowPriorityIO</key>
	<true/>
	<key>Nice</key>
	<integer>20</integer>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/mpPlanB</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>7207</integer>
	<key>TimeOut</key>
	<integer>1800</integer>
	<key>UserName</key>
	<string>root</string>
	<key>WorkingDirectory</key>
	<string>/</string>
</dict>
</plist>
