# Packages/Client/Files/Library/LaunchAgents/gov.llnl.MPLoginAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>gov.llnl.MPLoginAgent</string>
	<key>LimitLoadToSessionType</key>
	<string>LoginWindow</string>
	<key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/private/tmp/.MPAuthRun</key>
            <true/>
        </dict>
    </dict>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/PrivilegedHelperTools/MPLoginAgent.app/Contents/MacOS/MPLoginAgent</string>
	</array>
</dict>
</plist>
